# conf/templates/port-forward-ingress.service
[Unit]
Description=kubectl port-forward svc ingress
After=network.target

[Service]
Type=simple
ExecStart={{ script }}
Restart={{ restart_on | default('on-failure') }}
RestartSec=5s

[Install]
WantedBy=multi-user.target